# deploy/systemd/smart-mirror-usb-power.path
[Unit]
Description=Smart Mirror USB power state watcher
After=docker.service

[Path]
PathModified=/home/smartmirror/Downloads/smart-mirror/backend/data/usb-power-state
Unit=smart-mirror-usb-power.service

[Install]
WantedBy=multi-user.target
